# bluetooth-pan.service
[Unit]
Description=Bluetooth PAN Network Access Point
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple

# remove pan0 if it already exists, then recreate clean
ExecStartPre=/bin/bash -c 'ip link set pan0 down 2>/dev/null; ip link delete pan0 2>/dev/null; true'
ExecStartPre=/usr/bin/ip link add name pan0 type bridge
# disable STP — without this, new bridge ports wait ~30 s before forwarding,
# causing Android's DHCP client to time out before bnep0 is ready
ExecStartPre=/usr/bin/ip link set pan0 type bridge stp_state 0
ExecStartPre=/usr/bin/ip link set pan0 up
ExecStartPre=/usr/bin/ip addr add 192.168.50.1/24 dev pan0
# ULA IPv6 address — required for dnsmasq to send Router Advertisements
# so Android BT PAN clients receive an RA and configure themselves via SLAAC
ExecStartPre=/usr/sbin/ip -6 addr add fd00::1/64 dev pan0
# if br_netfilter is loaded, bridge frames go through iptables FORWARD chain;
# disabling it avoids DHCP/HTTP being dropped by a default DROP policy
ExecStartPre=/bin/bash -c 'sysctl -w net.bridge.bridge-nf-call-iptables=0 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'sysctl -w net.bridge.bridge-nf-call-ip6tables=0 2>/dev/null || true'
# Docker sets FORWARD policy to DROP and appends a catch-all DROP rule.
# Insert ACCEPT rules at the top of the chain so pan0 traffic is always allowed,
# regardless of Docker's rules or br_netfilter state.
ExecStartPre=/bin/bash -c 'ip6tables -C INPUT  -i pan0 -j ACCEPT 2>/dev/null || ip6tables -I INPUT  -i pan0 -j ACCEPT'
ExecStartPre=/bin/bash -c 'ip6tables -C FORWARD -i pan0 -j ACCEPT 2>/dev/null || ip6tables -I FORWARD -i pan0 -j ACCEPT'
ExecStartPre=/bin/bash -c 'iptables  -C INPUT  -i pan0 -j ACCEPT 2>/dev/null || iptables  -I INPUT  -i pan0 -j ACCEPT'
ExecStartPre=/bin/bash -c 'iptables  -C FORWARD -i pan0 -j ACCEPT 2>/dev/null || iptables  -I FORWARD -i pan0 -j ACCEPT'

# make Pi discoverable and pairable
ExecStartPre=/usr/bin/bluetoothctl power on
ExecStartPre=/usr/bin/bluetoothctl discoverable on
ExecStartPre=/usr/bin/bluetoothctl pairable on

# main process — Bluetooth NAP server
ExecStart=/usr/bin/bt-network -s nap pan0

# restart dnsmasq AFTER pan0 is up so it binds to the interface correctly
ExecStartPost=/bin/systemctl restart dnsmasq

# cleanup on stop
ExecStop=/bin/systemctl stop dnsmasq
ExecStop=/usr/bin/bluetoothctl discoverable off
ExecStop=/bin/bash -c 'ip6tables -D INPUT  -i pan0 -j ACCEPT 2>/dev/null; true'
ExecStop=/bin/bash -c 'ip6tables -D FORWARD -i pan0 -j ACCEPT 2>/dev/null; true'
ExecStop=/bin/bash -c 'iptables  -D INPUT  -i pan0 -j ACCEPT 2>/dev/null; true'
ExecStop=/bin/bash -c 'iptables  -D FORWARD -i pan0 -j ACCEPT 2>/dev/null; true'
ExecStop=/bin/bash -c 'ip link set pan0 down && ip link delete pan0 type bridge || true'

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
